# prod/pips/src/Scripts/dev/fabien_verone.cron
# $Id: fabien_verone.cron 23254 2016-11-01 07:16:31Z coelho $
#
# Debian Jessie i686
#
# compilation & validation once a day
25 2 * * * /home/fabien/32BIT/pips_check_compile --force --report --source /home/fabien/32BIT/pipsrc.sh coelho@cri.ensmp.fr /home/fabien/32BIT/prod /home/fabien/32BIT/validation
