# prod/pips/src/Scripts/dev/fabien_compiegne.cron
# $Id: fabien_compiegne.cron 23209 2016-10-18 08:49:36Z coelho $
#
# Debian jessie x86_64
#
# CHECK: "on commit" runs on every minute
#
* * * * * /home/fabien/CHECK/pips_check_compile --source /home/fabien/CHECK/pipsrc.sh --keepfail --validation-option DO_SLOW= --lock=$HOME --lock-delay=7200 pips-validation@cri.ensmp.fr /home/fabien/CHECK/prod /home/fabien/CHECK/validation
#
# & force once a day
#
7 0 * * * /home/fabien/CHECK/pips_check_compile --wait --force --report --source /home/fabien/CHECK/pipsrc.sh --validation-option DO_SLOW= --lock=$HOME --lock-delay=7200 pips-validation@cri.ensmp.fr /home/fabien/CHECK/prod /home/fabien/CHECK/validation
#
# AUTOC: test autoconf version on every hour
#
@hourly /home/fabien/AUTOC/pips_check_compile --autoconf --autoconf-enable='pyps hpfc gpips fortran95' --keepfail --source=/home/fabien/AUTOC/pipsrc.sh --wait --lock=$HOME --lock-delay=7200 pips-validation@cri.ensmp.fr /home/fabien/AUTOC/prod /home/fabien/AUTOC/validation
#
# & force once a day
#
7 1 * * * /home/fabien/AUTOC/pips_check_compile --wait --force --report --autoconf --autoconf-enable='pyps hpfc gpips fortran95' --keepfail --source=/home/fabien/AUTOC/pipsrc.sh --validation-option PIPS_VALIDATION_EXE=1 --lock=$HOME --lock-delay=7200 pips-validation@cri.ensmp.fr /home/fabien/AUTOC/prod /home/fabien/AUTOC/validation
#
# CROSS: 32bits cross compilation, without validation
#
7 2 * * * /home/fabien/CROSS/pips_check_compile --force --report --source=/home/fabien/CROSS/pipsrc.sh pips-validation@cri.ensmp.fr /home/fabien/CROSS/prod
#
# PRIVE: run private validation
#
7 3 * * * /home/fabien/PRIVE/pips_check_compile --force --report --source /home/fabien/PRIVE/pipsrc.sh pips-private@cri.ensmp.fr /home/fabien/PRIVE/prod /home/fabien/PRIVE/private
#
# DEB80: debian jessie *with* execution
#
7 6 * * * /home/fabien/DEB80/pips_check_compile --force --report --source /home/fabien/DEB80/pipsrc.sh coelho@cri.ensmp.fr /home/fabien/DEB80/prod /home/fabien/DEB80/validation
