# prod/pips/src/Scripts/dev/fabien_pau.cron
#
# $Id: fabien_pau.cron 23134 2016-07-07 05:49:54Z coelho $
#
# Ubuntu 16.04 Xenial x86_64 on PAU
#
# autoconf force once a day
#
7 1 * * * /home/fabien/AUTOP/pips_check_compile --wait --force --report --autoconf --keepfail --source=/home/fabien/AUTOP/pipsrc.sh --validation-option PIPS_VALIDATION_EXE=1 pips-validation@cri.ensmp.fr /home/fabien/AUTOP/prod /home/fabien/AUTOP/validation
#
# and try once an hour
#
37 * * * * /home/fabien/AUTOP/pips_check_compile --autoconf --keepfail --source=/home/fabien/AUTOP/pipsrc.sh --validation-option PIPS_VALIDATION_EXE=1 pips-validation@cri.ensmp.fr /home/fabien/AUTOP/prod /home/fabien/AUTOP/validation
#
# clang once a day
#
7 2 * * * /home/fabien/CLANG/pips_check_compile --wait --force --report --keepfail --source=/home/fabien/CLANG/pipsrc.sh pips-validation@cri.ensmp.fr /home/fabien/CLANG/prod /home/fabien/CLANG/validation
#
# clang trunk (3.8.0)
#
37 2 * * * /home/fabien/CLGTK/pips_check_compile --wait --force --report --keepfail --source=/home/fabien/CLGTK/pipsrc.sh pips-validation@cri.ensmp.fr /home/fabien/CLGTK/prod /home/fabien/CLGTK/validation
#
# gcc5x once a day
#
37 3 * * * /home/fabien/GCC5X/pips_check_compile --wait --force --report --keepfail --source=/home/fabien/GCC5X/pipsrc.sh pips-validation@cri.ensmp.fr /home/fabien/GCC5X/prod /home/fabien/GCC5X/validation
#
# gcc6x once a day
#
7 4 * * * /home/fabien/GCC6X/pips_check_compile --wait --force --report --keepfail --source=/home/fabien/GCC6X/pipsrc.sh pips-validation@cri.ensmp.fr /home/fabien/GCC6X/prod /home/fabien/GCC6X/validation
#
# gcctk once a day
#
7 5 * * * /home/fabien/GCCTK/pips_check_compile --wait --force --report --keepfail --source=/home/fabien/GCCTK/pipsrc.sh pips-validation@cri.ensmp.fr /home/fabien/GCCTK/prod /home/fabien/GCCTK/validation
#
# spear/cmini/fmini once an hour + force once a day
#
5 4 * * * /home/fabien/SPEAR/pips_check_compile --wait --force --report --keepfail --source=/home/fabien/SPEAR/pipsrc.sh --validation-option TARGET=Spear pips-validation@cri.ensmp.fr /home/fabien/SPEAR/prod /home/fabien/SPEAR/pipspriv
15 4 * * * /home/fabien/CMINI/pips_check_compile --wait --force --report --keepfail --source=/home/fabien/CMINI/pipsrc.sh --validation-option TARGET=Cmini pips-validation@cri.ensmp.fr /home/fabien/CMINI/prod /home/fabien/CMINI/validation
25 4 * * * /home/fabien/FMINI/pips_check_compile --wait --force --report --keepfail --source=/home/fabien/FMINI/pipsrc.sh --validation-option TARGET=Fmini pips-validation@cri.ensmp.fr /home/fabien/FMINI/prod /home/fabien/FMINI/validation
#
35 * * * * /home/fabien/SPEAR/pips_check_compile --keepfail --source=/home/fabien/SPEAR/pipsrc.sh --validation-option TARGET=Spear pips-validation@cri.ensmp.fr /home/fabien/SPEAR/prod /home/fabien/SPEAR/pipspriv
45 * * * * /home/fabien/CMINI/pips_check_compile --keepfail --source=/home/fabien/CMINI/pipsrc.sh --validation-option TARGET=Cmini pips-validation@cri.ensmp.fr /home/fabien/CMINI/prod /home/fabien/CMINI/validation
55 * * * * /home/fabien/FMINI/pips_check_compile --keepfail --source=/home/fabien/FMINI/pipsrc.sh --validation-option TARGET=Fmini pips-validation@cri.ensmp.fr /home/fabien/FMINI/prod /home/fabien/FMINI/validation
#
# on commit check on pau, test every minute and force once a day
#
* * * * * /home/fabien/PAUCK/pips_check_compile --source /home/fabien/PAUCK/pipsrc.sh --keepfail --validation-option DO_SLOW= pips-validation@cri.ensmp.fr /home/fabien/PAUCK/prod /home/fabien/PAUCK/validation
7 0 * * * /home/fabien/PAUCK/pips_check_compile --wait --force --report --source /home/fabien/PAUCK/pipsrc.sh --keepfail --validation-option DO_SLOW= pips-validation@cri.ensmp.fr /home/fabien/PAUCK/prod /home/fabien/PAUCK/validation
